# services/core-node/src/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock

[supervisord]
nodaemon=true
user=root

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:tor]
command=/bin/bash -c "chown -R root:root /var/lib/tor 2>/dev/null; exec /usr/local/bin/tor-start.sh"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# VNC server on port 5907 — only starts when VNC_PASSWORD is set
# (otherwise /start-vnc.sh parks). Connect via SSH tunnel: vnc://localhost:5907
[program:vnc]
command=/start-vnc.sh
autostart=true
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=5
startsecs=2

# Core node - runs bitcoind (headless) or bitcoin-qt (GUI) based on GUI_MODE env var
# Set GUI_MODE=true to run with Qt GUI, otherwise runs headless bitcoind
[program:node]
command=/build/bcore/start_node.sh %(ENV_NODE_START_FLAGS)s
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=GUI_MODE="%(ENV_GUI_MODE)s"
priority=10
startsecs=5

[program:mining]
command=./start_mining.sh
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=WALLET_ADDRESS="%(ENV_WALLET_ADDRESS)s",MINING_CHECK_INTERVAL="%(ENV_CHECK_INTERVAL)s"
priority=30
startsecs=10
startretries=10

[program:api_server]
command=python3 /app/api_server.py
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=API_KEY="%(ENV_MODEL_API_KEY)s"
priority=20
startsecs=30
startretries=5

[program:bootstrap_peers]
command=/usr/local/bin/bootstrap-peers.sh
autostart=true
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=12
startsecs=0

[program:vanity_onion]
command=/usr/local/bin/vanity-onion-gen.sh
autostart=%(ENV_VANITY_ONION_ENABLED)s
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=15
startsecs=10

# Consumer alternative to vanity_onion: claim a pre-ground onion (e.g. `tensorc`)
# from an operator-provided onion-grinder pool + install via SIGHUP (no Tor
# cold-boot). Mutually exclusive with vanity_onion — a claim node sets
# ONION_CLAIM_ENABLED=true and VANITY_ONION_ENABLED=false. Default off (ENV
# default in the Dockerfile).
[program:onion_claim]
command=python3 /usr/local/bin/onion-claim-agent.py
autostart=%(ENV_ONION_CLAIM_ENABLED)s
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=15
startsecs=5

# Minimal Tor/peer observability (report only, no auto-remediation)
[program:tor_health]
command=/usr/local/bin/tor-health.sh
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=40
startsecs=0
